# agent/cyberx-agent.service
[Unit]
Description=CyberX Instance Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/cyberx-agent/cyberx-agent.py
Restart=on-failure
RestartSec=10
EnvironmentFile=-/etc/cyberx-agent/agent.env

# Security hardening
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/etc/wireguard

[Install]
WantedBy=multi-user.target
